# daemons/etc/service/temper_daemon/run
#!/bin/sh
exec /usr/local/sbin/temper_daemon.sh
